# deploy/systemd/teleuserbot.service
[Unit]
Description=TeleUserBot News Intelligence Userbot
After=network.target
Wants=network.target

[Service]
User=therayyanawaz
Group=therayyanawaz
WorkingDirectory=/home/therayyanawaz/Documents/Github/TeleUserBot-main
EnvironmentFile=/home/therayyanawaz/Documents/Github/TeleUserBot-main/.env
ExecStart=/home/therayyanawaz/Documents/Github/TeleUserBot-main/.venv/bin/python3 /home/therayyanawaz/Documents/Github/TeleUserBot-main/main.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

# Secure the process
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
